# files/etc/supervisor/conf.d/memcached.conf
[program:memcached]
command = /usr/bin/memcached -u root
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
